# FanController.service.in
[Unit]
Description = controlls fans connected to GPIO pins

[Service]
Type = simple
User = root
ExecStart  = @Program_Path@
ExecStop   = /bin/kill -s SIGINT $MAINPID
ExecReload = /bin/kill -s SIGHUP $MAINPID
Restart    = on-failure
KillSignal  = SIGINT
SendSIGKILL = false

[Install]
WantedBy = multi-user.target
